# deploy/systemd/polymarket-bot.service
[Unit]
Description=Polymarket Paper Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__PROJECT_DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/venv/bin/python -u src/main.py
Restart=always
RestartSec=10
StandardOutput=append:__PROJECT_DIR__/logs/bot.log
StandardError=append:__PROJECT_DIR__/logs/bot.error.log

[Install]
WantedBy=multi-user.target
